# crontab
*/3 * * * * root . /env.sh; /wu-pws.sh > /tmp/output.log
